# cloudflare-ddns.service
[Unit]
Description=Update DNS entry for this host to current IP

[Service]
Type=oneshot
ExecStart=/usr/bin/cloudflare-ddns
