# cloudflare-ddns.timer
[Unit]
Description=Update DNS entry in cloudflare

[Timer]
OnBootSec=1min
OnCalendar=*:0/2
Unit=cloudflare-ddns.service

[Install]
WantedBy=basic.target
